# bin/edbc.rc.sh.in
#!/bin/bash
#
# chkconfig: 345 99 01
# description: Icinga network monitor
### BEGIN INIT INFO
# Provides:          edbc
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: EventDB correlator Initscript
# Description:       EventDB Agent with aggregation and correlation
### END INIT INFO

. /etc/init.d/functions

USER=icinga
GROUP=icingacmd
LOCKFILE=@prefix@/var/edbc.lock
DAEMON=@prefix@/lib/edbc.py
PROG=edbc.py

RETVAL=0


start() {
        [ -x $DAEMON ] || exit 5

                        echo -n $"Starting $PROG: "
                        sudo -u $USER "$DAEMON" &
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && touch $LOCKFILE
        return $RETVAL
        }


stop() {
       echo -n $"Stopping $PROG:"
       killproc $DAEMON
       RETVAL=$?
       echo
       [ $RETVAL -eq 0 ] && rm -f $LOCKFILE
       return $RETVAL
       }

restart() {
   stop
   start
}

case "$1" in
start)
        status $PROG && exit 0
        $1
        ;;
stop)
        status $PROG
                        [ ! $? -eq 0 ] && exit 0
                        $1
        ;;
status)
        status edbc.py
        ;;
restart)
        $1
        ;;
*)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 2
        ;;
esac
exit $?
